# deployment-files/ha/fleet-follows-primary.service.in
[Unit]
Description=Start Proto Fleet only on the local HA database primary
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_ROOT__
ExecStart=__PROJECT_ROOT__/ha/fleet-follows-primary.sh
